# sausage-store-infra/ansible/roles/backend/templates/backend.service.j2
[Unit]
Description=Backend Service
After=network.target

[Service]
User=backend
WorkingDirectory=/opt/backend
ExecStart=/usr/bin/java -jar /opt/backend/backend.jar
Restart=always

[Install]
WantedBy=multi-user.target
